# lumina.service
[Unit]
Description=Lumina — Inspiration Manager (Next.js)
After=network.target

[Service]
Type=simple
User=uri
WorkingDirectory=/usr/local/lumina
EnvironmentFile=-/usr/local/lumina/.env.local
Environment=NODE_ENV=production
Environment=PORT=3009
ExecStart=/usr/bin/node server.compiled.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
